# build/package/webapplication/supervisor/conf.d/webapplication.conf
[program:webapplication]
directory=/app
command=/app/webapplication
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
autostart=false
autorestart=true
stopsignal=TERM
numprocs=1
